# change_wallpaper.service
[Unit]
Description=Change Wallpaper

[Service]
ExecStart=/home/$USER/change_wallpaper.sh

